# deploy_tools/gunicorn-upstart.template.conf
start on net-device-up
stop on shutdown

respawn

setuid alistair

chdir /home/alistair/sites/SITENAME/source

exec ../virtualenv/bin/gunicorn \
    --bind unix:/tmp/SITENAME.socket \
    --access-logfile ../access.log \
    --error-logfile ../error.log \
    superlists.wsgi:application
